# xdg-desktop-portal.service
[Unit]
Description=Portal service (xdg-desktop-portal)
# Depends on D-Bus only — graphical-session.target is NOT required here.
# KWin/plasmalogin activate the session without that target.
After=dbus.socket
Requires=xdg-desktop-portal-kde.service

[Service]
Type=dbus
BusName=org.freedesktop.portal.Desktop
ExecStart=/usr/libexec/xdg-desktop-portal
Restart=on-failure
RestartSec=3
Slice=session.slice

[Install]
WantedBy=default.target
